# deploy/bot.service
[Unit]
# Canli bahis takip ve Telegram bildirim botu icin systemd unit dosyasi.
# Kurulum:
#   sudo cp deploy/bot.service /etc/systemd/system/betting-bot.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now betting-bot
#   sudo journalctl -u betting-bot -f       # log takibi
#
# Asagidaki User ve dizin yollari Oracle Cloud Ubuntu varsayilanlarina gore
# ayarlandi. Repo'yu farkli bir yere klonladiysan WorkingDirectory ve
# ExecStart yollarini buna gore guncelle.

Description=Canli bahis takip ve Telegram bildirim botu
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/claude-test/python
ExecStart=/home/ubuntu/claude-test/.venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
